# contrib/debian/preremove.sh
#!/bin/bash
systemctl stop haos-agent
systemctl disable haos-agent
